# packaging/systemd/agent-forkd.service
[Unit]
Description=Forked development environment manager
After=local-fs.target systemd-networkd.service systemd-machined.service
Requires=systemd-networkd.service systemd-machined.service

[Service]
Type=simple
Environment=AGENTFS=/agentfs
Environment=AGENT_FORKD_CONFIG=/etc/agent-forkd/config.json
ExecStart=/usr/local/bin/agent-forkd --agentfs ${AGENTFS} --config ${AGENT_FORKD_CONFIG}
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=multi-user.target
